# harbi-basketball-events.service
[Unit]
Description=Harbi Basketball Arbitrage Event Loop
After=network.target

[Service]
Type=simple
User=giray
WorkingDirectory=/home/giray/harbi
ExecStart=/home/giray/harbi/venv/bin/python3 -u /home/giray/harbi/run_basketball_event_loop.py
Restart=always
RestartSec=10
StandardOutput=append:/home/giray/harbi/logs/basketball_event_loop.log
StandardError=append:/home/giray/harbi/logs/basketball_event_loop.log

[Install]
WantedBy=multi-user.target
